# puppet/modules/cron/files/debomatic.cron
MAILTO=and@debian.org
*/5 * * * * and		cd $HOME/public_html && { rsync -avz --recursive --delete --links --exclude "unstable/aptcache" --exclude "unstable/build" --exclude "unstable/gpg" --exclude "unstable/unstable" -e "ssh -q -i $HOME/.ssh/rsync" mrlawrence:debomatic . ; } > /dev/null
